# files/services/network_setup.service
[Unit]
Description=Network setting service
Before=network.target

[Service]
Environment=USERNAME=robot
ExecStart=/home/robot/setup_scripts/setup_network.sh
Type=oneshot
User=root

[Install]
WantedBy=multi-user.target
